# systemd/coffee-updater.timer
[Unit]
Description=Coffee Manager Auto-Updater Timer
Requires=coffee-updater.service

[Timer]
# Run 10 minutes after boot
OnBootSec=10min
# Run daily at 2:00 AM
OnCalendar=daily 02:00
# Catch up on missed runs if system was off
Persistent=true

[Install]
WantedBy=timers.target


